# oci/vm/config/supervisord.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:sway]
command=/opt/browserbird/runtime/sway_run.sh
user=bbuser
environment=HOME="/home/bbuser"
autostart=true
autorestart=true
stderr_logfile=/var/log/sway.err.log
stdout_logfile=/var/log/sway.out.log

[program:wayvnc]
command=/opt/browserbird/runtime/wayvnc_run.sh
user=bbuser
environment=HOME="/home/bbuser"
autostart=true
autorestart=true
stderr_logfile=/var/log/wayvnc.err.log
stdout_logfile=/var/log/wayvnc.out.log

[program:novnc]
command=/opt/browserbird/runtime/novnc_run.sh
user=bbuser
environment=HOME="/home/bbuser"
autostart=true
autorestart=true
stderr_logfile=/var/log/novnc.err.log
stdout_logfile=/var/log/novnc.out.log

[program:playwright-mcp]
command=/opt/browserbird/runtime/playwright_mcp_run.sh
user=bbuser
environment=HOME="/home/bbuser"
autostart=true
autorestart=true
startretries=10
stderr_logfile=/var/log/playwright-mcp.err.log
stdout_logfile=/var/log/playwright-mcp.out.log

[program:file-server]
command=node /opt/browserbird/runtime/file_server.ts
user=bbuser
environment=HOME="/home/bbuser"
autostart=true
autorestart=true
stderr_logfile=/var/log/file-server.err.log
stdout_logfile=/var/log/file-server.out.log
